# debian/pywebdriver.conf
[program:pywebdriver]
user=pywebdriver
command=/usr/bin/pywebdriverd
stdout_logfile=/var/log/supervisor/pywebdriver.log
redirect_stderr=true
autostart=true
autorestart=true
